# crontab
*/30 * * * * cd /code && python3 main.py >> /logfile
